# crontab
crontab -e
@reboot /usr/bin/python /home/pi/pushover/Push_Restart.py
